# screenshot-scheduler/systemd/screenshot_scheduler.service
[Unit]
Description=Webcam Screenshot Scheduler
After=network.target

[Service]
Type=simple
ExecStart=/path/to/screenshot_scheduler.py
Restart=always

[Install]
WantedBy=multi-user.target
